# tools/cloud/pulsebrief-web.service
[Unit]
Description=Pulse Brief web server
After=network-online.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=pulsebrief
Group=pulsebrief
WorkingDirectory=/opt/pulsebrief/web
EnvironmentFile=-/etc/pulsebrief/pulsebrief.env
ExecStart=/usr/bin/dotnet /opt/pulsebrief/web/PulseBrief.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=pulsebrief-web
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
