# tools/cloud/pulsebrief-mongodb-backup.service
[Unit]
Description=Pulse Brief MongoDB backup
After=mongod.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/pulsebrief-mongodb-backup.sh
